# ci/docker/clear-log.cron
59 23 * * * rm -f /var/log/nginx/access.log /var/log/nginx/error.log && kill -USR1 $(cat /usr/local/openresty/nginx/logs/nginx.pid)